# config/systemd/isaac-robot.service
[Unit]
Description=Isaac Robot System
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=nano
Group=nano
Environment="ISAAC_DEV_DIR=/home/nano/src/jetson-orin-nano"
Environment="ISAAC_INSTALL_DIR=/opt/isaac-robot"
Environment="PYTHONPATH=/home/nano/src/jetson-orin-nano:/opt/isaac-robot"

# Find Isaac root (dev takes precedence) and start robot
ExecStart=/bin/bash -c 'ISAAC_ROOT=$(python3 /home/nano/src/jetson-orin-nano/scripts/utils/find_isaac_root.py 2>/dev/null || python3 /opt/isaac-robot/scripts/utils/find_isaac_root.py 2>/dev/null || echo "/opt/isaac-robot"); cd "$ISAAC_ROOT"; exec "$ISAAC_ROOT/scripts/system/start_robot.sh"'

# Restart on failure with exponential backoff
Restart=on-failure
RestartSec=10
RestartSteps=5
RestartMaxDelaySec=300

# Watchdog
WatchdogSec=60

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
