# config/systemd/isaac-cleanup.timer
[Unit]
Description=Isaac Robot System Cleanup Timer
After=network-online.target

[Timer]
# Run weekly on Sunday at 2 AM
OnCalendar=weekly
OnCalendar=Sun 02:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
